# recipes-multimedia/imdt-camera-demo/files/system/imdt-video-stream-2.service
#
# Copyright (c) 2021 IMD Technologies
#

[Unit]
Description=Video Stream from CSI2

[Service]
ExecStart=/opt/imdt/camera-demo/bin/camera2-stream.sh

